# systemd/mosaic-usb0.service
[Unit]
Description=Configure static IP for usb0 (RPi side, Mosaic-x5 USB-Ethernet hybrid)
After=NetworkManager.service network-online.target
Requires=NetworkManager.service
Wants=network-online.target
ConditionPathExists=/sys/class/net/usb0

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure_usb0_nmcli.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

